# scripts/develed.service
[Unit]
Description=Develed Bot

[Service]
Type=simple
ExecStart=/usr/bin/bot

[Install]
WantedBy=multi-user.target
