# scripts/dspd.service
[Unit]
Description=Develed dspd daemon

[Service]
Type=simple
ExecStart=/usr/bin/dspd

[Install]
WantedBy=multi-user.target
